# roles/streamwall/templates/vnc.socket.j2
[Unit]
Description=VNC Socket
PartOf=vnc.service

[Socket]
ListenStream=127.0.0.1:5900
Accept=true

[Install]
WantedBy=sockets.target
