# roles/streambackup/files/streambackup.service.j2
[Unit]
Description=Stream VOD Backup Service
PartOf=streambackup.target

[Service]
Restart=always
RestartSec=1s
ExecStart=/home/streamwall/streambackup.sh

[Install]
WantedBy=streambackup.target
