# roles/streambackup/files/streambackup.target
[Unit]
Description=Streambackup Service
BindsTo=stream.target
After=stream.target

[Install]
WantedBy=stream.target
