# roles/streamwall/files/stream.target
[Unit]
Description=Streamwall Stream
BindsTo=graphical-session.target
After=graphical-session.target

[Install]
WantedBy=graphical-session.target
